# deploy/password-manager.service
[Unit]
Description=Password Manager Web Application
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=passwordmgr
Group=passwordmgr
WorkingDirectory=/opt/password-manager
ExecStart=/opt/password-manager/password-manager
Environment=PORT=8081
Environment=GIN_MODE=release

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/password-manager

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=password-manager

[Install]
WantedBy=multi-user.target
